# deploy/ams-backend.service
[Unit]
Description=AMS Config Manager Backend
After=network.target

[Service]
Type=simple
User=actions-runner
Group=actions-runner
WorkingDirectory=/opt/ams/backend
ExecStart=/opt/ams/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/ams/backend

[Install]
WantedBy=multi-user.target
